# conf/django_blog.ini
[program:django_blog]
command=/root/django_blog/env/bin/gunicorn blog.wsgi -c /root/django_blog/gun.py
directory=/root/django_blog
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile =/root/django_blog/blog_stdout.log  
stderr_logfile =/root/django_blog/blog_stderr.log 
stdout_logfile_maxbytes =500MB  
stdout_logfile_backups =50  
stdout_capture_maxbytes =1MB  
stdout_events_enabled =false
environment=PYTHONPATH=/root/django_blog:/root/django_blog/django-comments-xtd:/root/django_blog/zinnia-theme-bootstrap,DJANGO_SETTINGS_MODULE=blog.settings_production,DGDB_PASSWD='your mysql user passwd',DGDB_USER='your mysql user name',SECRET_KEY='your secret key'
